# infra/systemd/vecta-firstboot.service
[Unit]
Description=Vecta KMS First-Boot Wizard
After=network-online.target docker.service
Wants=network-online.target
ConditionPathExists=!/etc/vecta/deployment.yaml

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
Group=root
WorkingDirectory=/opt/vecta
Environment=FIRSTBOOT_CONFIG_DIR=/etc/vecta
ExecStartPre=/usr/bin/mkdir -p /etc/vecta
ExecStart=/bin/bash /opt/vecta/infra/scripts/compose-kms.sh --profile firstboot up -d firstboot
ExecStop=/bin/bash /opt/vecta/infra/scripts/compose-kms.sh stop firstboot
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
