# infra/systemd/vecta-containers.target
[Unit]
Description=Vecta KMS Individual Container Units
Requires=docker.service
After=docker.service
