# infra/systemd/vecta-deployment.path
[Unit]
Description=Watch Vecta deployment config and trigger stack startup

[Path]
PathExists=/etc/vecta/deployment.yaml
PathChanged=/etc/vecta/deployment.yaml
PathModified=/etc/vecta/deployment.yaml
Unit=vecta-stack.service

[Install]
WantedBy=multi-user.target
