# infra/systemd/containers/vecta-certs.service
[Unit]
Description=Vecta KMS Container (certs)
After=docker.service network-online.target
Wants=network-online.target
PartOf=vecta-containers.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/vecta
EnvironmentFile=-/etc/vecta/vecta.env
ExecStart=/usr/bin/docker compose -f /opt/vecta/docker-compose.yml up -d certs
ExecStop=/usr/bin/docker compose -f /opt/vecta/docker-compose.yml stop certs
ExecReload=/usr/bin/docker compose -f /opt/vecta/docker-compose.yml up -d certs
TimeoutStartSec=300
TimeoutStopSec=120

[Install]
WantedBy=vecta-containers.target