# infra/systemd/containers/vecta-consul.service
[Unit]
Description=Vecta KMS Container (consul)
After=docker.service network-online.target
Wants=network-online.target
PartOf=vecta-containers.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/vecta
EnvironmentFile=-/etc/vecta/vecta.env
ExecStart=/bin/bash /opt/vecta/infra/scripts/compose-kms.sh up -d consul
ExecStop=/bin/bash /opt/vecta/infra/scripts/compose-kms.sh stop consul
ExecReload=/bin/bash /opt/vecta/infra/scripts/compose-kms.sh up -d consul
TimeoutStartSec=300
TimeoutStopSec=120

[Install]
WantedBy=vecta-containers.target
