# infra/systemd/containers/vecta-discovery.service
[Unit]
Description=Vecta KMS Container (discovery)
After=docker.service network-online.target
Wants=network-online.target
PartOf=vecta-containers.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/vecta
EnvironmentFile=-/etc/vecta/vecta.env
ExecStart=/usr/bin/docker compose -f /opt/vecta/docker-compose.yml up -d discovery
ExecStop=/usr/bin/docker compose -f /opt/vecta/docker-compose.yml stop discovery
ExecReload=/usr/bin/docker compose -f /opt/vecta/docker-compose.yml up -d discovery
TimeoutStartSec=300
TimeoutStopSec=120

[Install]
WantedBy=vecta-containers.target