# infra/systemd/containers/vecta-reporting.service
[Unit]
Description=Vecta KMS Container (reporting)
After=docker.service network-online.target
Wants=network-online.target
PartOf=vecta-containers.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/vecta
EnvironmentFile=-/etc/vecta/vecta.env
ExecStart=/bin/bash /opt/vecta/infra/scripts/compose-kms.sh up -d reporting
ExecStop=/bin/bash /opt/vecta/infra/scripts/compose-kms.sh stop reporting
ExecReload=/bin/bash /opt/vecta/infra/scripts/compose-kms.sh up -d reporting
TimeoutStartSec=300
TimeoutStopSec=120

[Install]
WantedBy=vecta-containers.target
