# infra/systemd/containers/vecta-compliance.service
[Unit]
Description=Vecta KMS Container (compliance)
After=docker.service network-online.target
Wants=network-online.target
PartOf=vecta-containers.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/vecta
EnvironmentFile=-/etc/vecta/vecta.env
ExecStart=/bin/bash /opt/vecta/infra/scripts/compose-kms.sh up -d compliance
ExecStop=/bin/bash /opt/vecta/infra/scripts/compose-kms.sh stop compliance
ExecReload=/bin/bash /opt/vecta/infra/scripts/compose-kms.sh up -d compliance
TimeoutStartSec=300
TimeoutStopSec=120

[Install]
WantedBy=vecta-containers.target
