# infra/systemd/containers/vecta-certs.service
[Unit]
Description=Vecta KMS Container (certs)
After=docker.service network-online.target
Wants=network-online.target
PartOf=vecta-containers.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/vecta
EnvironmentFile=-/etc/vecta/vecta.env
ExecStart=/bin/bash /opt/vecta/infra/scripts/compose-kms.sh up -d certs
ExecStop=/bin/bash /opt/vecta/infra/scripts/compose-kms.sh stop certs
ExecReload=/bin/bash /opt/vecta/infra/scripts/compose-kms.sh up -d certs
TimeoutStartSec=300
TimeoutStopSec=120

[Install]
WantedBy=vecta-containers.target
