# infra/systemd/containers/vecta-cluster-manager.service
[Unit]
Description=Vecta KMS Container (cluster-manager)
After=docker.service network-online.target
Wants=network-online.target
PartOf=vecta-containers.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/vecta
EnvironmentFile=-/etc/vecta/vecta.env
ExecStart=/usr/bin/docker compose -f /opt/vecta/docker-compose.yml up -d cluster-manager
ExecStop=/usr/bin/docker compose -f /opt/vecta/docker-compose.yml stop cluster-manager
ExecReload=/usr/bin/docker compose -f /opt/vecta/docker-compose.yml up -d cluster-manager
TimeoutStartSec=300
TimeoutStopSec=120

[Install]
WantedBy=vecta-containers.target