# infra/systemd/containers/vecta-postgres.service
[Unit]
Description=Vecta KMS Container (postgres)
After=docker.service network-online.target
Wants=network-online.target
PartOf=vecta-containers.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/vecta
EnvironmentFile=-/etc/vecta/vecta.env
ExecStart=/usr/bin/docker compose -f /opt/vecta/docker-compose.yml up -d postgres
ExecStop=/usr/bin/docker compose -f /opt/vecta/docker-compose.yml stop postgres
ExecReload=/usr/bin/docker compose -f /opt/vecta/docker-compose.yml up -d postgres
TimeoutStartSec=300
TimeoutStopSec=120

[Install]
WantedBy=vecta-containers.target