# infra/systemd/containers/vecta-dataprotect.service
[Unit]
Description=Vecta KMS Container (dataprotect)
After=docker.service network-online.target
Wants=network-online.target
PartOf=vecta-containers.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/vecta
EnvironmentFile=-/etc/vecta/vecta.env
ExecStart=/bin/bash /opt/vecta/infra/scripts/compose-kms.sh up -d dataprotect
ExecStop=/bin/bash /opt/vecta/infra/scripts/compose-kms.sh stop dataprotect
ExecReload=/bin/bash /opt/vecta/infra/scripts/compose-kms.sh up -d dataprotect
TimeoutStartSec=300
TimeoutStopSec=120

[Install]
WantedBy=vecta-containers.target
